# mkosi.extra/usr/lib/udev/rules.d/99-incus-agent.rules
ACTION=="add", SYMLINK=="virtio-ports/org.linuxcontainers.incus", TAG+="systemd", ACTION=="add", RUN+="/bin/systemctl start incus-agent.service"
